# ansible/roles/alloy/templates/alloy.service.j2
[Unit]
Description=Grafana Alloy
Documentation=https://grafana.com/docs/alloy/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ alloy_user }}
Group={{ alloy_group }}
ExecStart={{ alloy_binary_path }} run {{ alloy_config_file }} \
  --storage.path={{ alloy_data_dir }} \
  --server.http.listen-addr=0.0.0.0:{{ alloy_http_port }}
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Log access
SupplementaryGroups=systemd-journal
ReadOnlyPaths=/var/log

# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
