# ansible/roles/alloy/templates/alloy.openrc.j2
#!/sbin/openrc-run
# Generated by Ansible on {{ ansible_date_time.iso8601 }}

description="Grafana Alloy - metrics and log collector"

command="{{ alloy_binary_path }}"
command_args="run {{ alloy_config_file }} --storage.path={{ alloy_data_dir }} --server.http.listen-addr=0.0.0.0:{{ alloy_http_port }}"
command_user="{{ alloy_user }}:{{ alloy_group }}"
command_background=true
directory="{{ alloy_data_dir }}"
pidfile="/run/alloy.pid"
output_log="{{ alloy_log_file }}"
error_log="{{ alloy_log_file }}"

depend() {
    need net
    after firewall
}

start_pre() {
    checkpath --directory --owner {{ alloy_user }}:{{ alloy_group }} --mode 0750 {{ alloy_data_dir }}
    checkpath --file --owner {{ alloy_user }}:{{ alloy_group }} --mode 0640 {{ alloy_log_file }}
}

reload() {
    ebegin "Reloading ${RC_SVCNAME}"
    start-stop-daemon --signal HUP --pidfile "${pidfile}"
    eend $?
}
